# misc/supervisord.conf
[supervisord]
nodaemon=true

[program:incron]
command=/usr/sbin/incrond -n -f /etc/incron.conf

[program:nginx]
command=nginx -g "daemon off;"
